# plugins/input-xflow/supervisord.conf
[supervisord]
nodaemon=true

[program:telegraf_xflow]
command=/bin/bash /data/templates/start-input-xflow.sh
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/telegraf-xflow.log
stdout_logfile=/var/log/telegraf-xflow.log

[program:telegraf_xflow_internal]
command=/bin/bash /data/templates/start-input-xflow_internal.sh
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/telegraf-xflow_internal.log
stdout_logfile=/var/log/telegraf-xflow_internal.log

[program:nfacctd]
command=/usr/sbin/nfacctd -f /data/templates/nfacctd.tmpl
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/nfacctd.log
stdout_logfile=/var/log/nfacctd.log

