# plugins/input-xflow/supervisord.conf
[supervisord]
nodaemon=true

[program:telegraf]
command=/data/templates/start-input-xflow.sh

[program:nfacctd]
command=/usr/sbin/nfacctd -f /data/templates/nfacctd.tmpl >>/var/log/nfacctd.log 2>&1

