# plugins/input-xflow/supervisord.conf
[supervisord]
nodaemon=true

[program:telegraf]
command=/data/templates/start-input-xflow.sh
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/telegraf.log
stdout_logfile=/var/log/telegraf.log


[program:nfacctd]
command=/usr/sbin/nfacctd -f /data/templates/nfacctd.tmpl
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/nfacctd.log
stdout_logfile=/var/log/nfacctd.log

